# cab/config/pi_supervisor/update.conf
[program:auto_update]
command = /usr/bin/python3 auto_downloader.py
directory = /home/pi/jprt/cab/shell/update/
process_name = auto_update
autostart = true
startsecs = 60
autorestart= true
;exitcodes = 0,2
user = pi
redirect_stderr = true
;stdout_logfile = /dev/null
stdout_logfile = /tmp/update.log
;environment = KEY1=value1,KEY2=value2
environment = HOME="/home/pi" 

